# build_files/systemd/flatpak-maintenance.timer
[Unit]
Description=Run daily Flatpak maintenance

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target